# support/daemontools/sshd/run
#!/bin/sh
exec /usr/sbin/sshd -p 2022 -D
